# root/etc/services.d/comlistener/run
#!/bin/bash

cd /comlistener || exit

exec \
  s6-setuidgid abc python3 /comlistener/comlistener.py
