# root/etc/services.d/comlistener/run
#!/bin/bash

cd /comlistener
./comlistener.py
